# overrides/etc/cron.hourly/0010-gc-dead-consul-servers
#!/bin/sh
yum install python-pip
python-pip install boto3 requests
exec python2.6 /usr/local/bin/gc-dead-consul-servers
